# EdgeProcessor/switch_poller.service
[Unit]
Description=IoT SOC - C9200-24T Switch Poller & Isolation Executor
Documentation=https://github.com/gavishka-sahan/Safenode-EdgeSystem
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Safenode-EdgeSystem
ExecStart=/usr/bin/python3 /home/pi/Safenode-EdgeSystem/switch_poller.py
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Logging — journald captures stdout/stderr
StandardOutput=journal
StandardError=journal
SyslogIdentifier=switch_poller

[Install]
WantedBy=multi-user.target
