# EdgeProcessor/switch_poller.service
[Unit]
Description=Safenode Switch Poller - Cisco C9200 Polling & VLAN Isolation
Documentation=https://github.com/gavishka-sahan/Safenode-EdgeSystem
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nomad
WorkingDirectory=/home/nomad/Safenode-EdgeSystem/EdgeProcessor
ExecStart=/usr/bin/python3 /home/nomad/Safenode-EdgeSystem/EdgeProcessor/switch_poller.py
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=SwitchPoller

[Install]
WantedBy=multi-user.target
