# EdgeProcessor/SwitchPoller.service
[Unit]
Description=Safenode Switch Poller - Cisco C3650 Polling & VLAN Isolation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/nomad/Safenode-EdgeSystem/EdgeProcessor
ExecStart=/usr/bin/python3 /home/nomad/Safenode-EdgeSystem/EdgeProcessor/SwitchPoller.py
Restart=always
RestartSec=10

Environment="PYTHONUNBUFFERED=1"

MemoryMax=512M
CPUQuota=100%

StandardOutput=append:/var/log/EdgeProcessor/SwitchPoller.log
StandardError=append:/var/log/EdgeProcessor/SwitchPoller.log

[Install]
WantedBy=multi-user.target
