# EdgeProcessor/ServiceFile/CloudAdapter.service
[Unit]
Description=IoT Network Feature Extractor
After=network.target
#Documentation=https://github.com/

[Service]
Type=simple
User=root
WorkingDirectory=/opt/CloudAdapter
ExecStart=/usr/bin/python3 /opt/CloudAdapter/CloudAdapter.py
Restart=always
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"

# Resource limits (Raspberry Pi 4GB)
MemoryMax=1G
CPUQuota=100%

# Logging
StandardOutput=append:/var/log/CloudAdapter/CloudAdapter.log
StandardError=append:/var/log/CloudAdapter/CloudAdapter.log

# Capabilities (for packet capture without full root)
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
