# EdgeProcessor/ServiceFile/FeatureLogReceiver.service
[Unit]
Description=IoT Network Edge Processor
After=network.target
#Documentation=https://github.com/

[Service]
Type=simple
User=root
WorkingDirectory=/opt/FeatureLogReceiver
ExecStart=/usr/bin/python3 /opt/FeatureLogReceiver/FeatureLogReceiver.py
Restart=always
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"

# Resource limits (Raspberry Pi 4GB)
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=append:/var/log/FeatureLogReceiver/FeatureLogReceiver.log
StandardError=append:/var/log/FeatureLogReceiver/FeatureLogReceiver.log

# Capabilities (for packet capture without full root)
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
